# systemd/user/spark-mcp.service.d/20-workspace-rust-mount.conf.in
[Service]
# Keep Rust workspace indexing enabled for local hover/search parity during mixed-language refactors.
Environment=SPARK_MCP_INCLUDE_WORKSPACE_RUST=1
