# systemd/user/spark-mcp.service.d/zz-no-build-on-restart.conf.in
[Service]
# Disable source-triggered builds on restart; run the existing release binary directly.
ExecStartPre=
