# systemd/user/spark-mcp.service.d/zz-restart-guardrails.conf.in
[Unit]
StartLimitIntervalSec=30min
StartLimitBurst=4

[Service]
Restart=on-failure
RestartSec=15s
RestartPreventExitStatus=127 203
